# social_home_dev/rootfs/etc/s6-overlay/s6-rc.d/socialhome/run
#!/command/with-contenv bashio
# ==============================================================================
# s6-overlay v3 longrun for the Social Home server (dev channel).
#
# ``/data/social_home.toml`` is rendered by the ``init-config``
# oneshot before this longrun starts (see
# ``dependencies.d/init-config``).
# ==============================================================================

bashio::log.info "Starting Social Home..."

exec python3 -m social_home.app --config /data/social_home.toml
