# social_home/rootfs/etc/s6-overlay/s6-rc.d/socialhome/run
#!/command/with-contenv bashio
# ==============================================================================
# s6-overlay v3 longrun for the Social Home server.
#
# ``/data/social_home.toml`` is rendered by the ``init-config``
# oneshot before this longrun starts (declared via
# ``dependencies.d/init-config``). Bootstrap logic (HA owner
# provisioning, integration token, Supervisor discovery push)
# lives in the Python core's ``HaBootstrap`` and fires on first
# start.
# ==============================================================================

bashio::log.info "Starting Social Home..."

# SUPERVISOR_TOKEN / SUPERVISOR_URL are injected by the HA
# Supervisor; the Python core picks them up from the environment.
exec python3 -m social_home.app --config /data/social_home.toml
